# debian/postinst
export ARCH=arm
export CROSS_COMPILE=/usr/bin/

python -c 'import Adafruit_BBIO.GPIO' || pip install --upgrade Adafruit_BBIO

systemctl daemon-reload
systemctl enable ospi.service
systemctl restart ospi.service
